# root/etc/services.d/subversion/run
#!/usr/bin/with-contenv sh

exec /usr/bin/svnserve -d --foreground -r /home/svn --listen-port 3690;
